# roles/ideascube/files/ideascube_leftover_files.service
[Unit]
Description=Run ideascube clean leftover-files
RefuseManualStart=no
RefuseManualStop=yes
 
[Service]
Type=oneshot
ExecStart=/usr/bin/ideascube clean leftover-files
